# backend/docker/frankenphp/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

# Serveur web FrankenPHP (Caddy + PHP en processus), drop vers www-data.
[program:frankenphp]
command=gosu www-data frankenphp run --config /etc/frankenphp/Caddyfile
autostart=true
autorestart=true
priority=10
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Worker Messenger : enrichissement Gemini, recherche de couvertures, push.
[program:messenger]
command=gosu www-data php /app/bin/console messenger:consume async --limit=500 --time-limit=3600
autostart=true
autorestart=true
priority=20
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Planificateur Symfony (transport scheduler).
[program:scheduler]
command=gosu www-data php /app/bin/console messenger:consume scheduler_default --time-limit=3600
autostart=true
autorestart=true
priority=20
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
